# content/cat-01-server-compute/UC-1.2.119.json
{
  "$schema": "../../schemas/uc.schema.json",
  "id": "1.2.119",
  "title": "Registry Run Key Persistence Monitoring",
  "criticality": "critical",
  "difficulty": "advanced",
  "monitoringType": [
    "Security"
  ],
  "splunkPillar": "Security",
  "dataSources": "`sourcetype=WinEventLog:Microsoft-Windows-Sysmon/Operational` (EventCode 13)",
  "app": "`Splunk_TA_windows`",
  "spl": "index=wineventlog EventCode=13\n| where match(TargetObject, \"(?i)(CurrentVersion\\\\\\\\Run|CurrentVersion\\\\\\\\RunOnce|Winlogon\\\\\\\\Shell|Winlogon\\\\\\\\Userinit|Explorer\\\\\\\\Shell Folders)\")\n| where NOT match(Details, \"(?i)(program files|windows\\\\\\\\system32|syswow64)\")\n| table _time, host, User, Image, TargetObject, Details\n| sort -_time",
  "description": "Registry Run keys are the most common persistence mechanism for malware. Monitoring autostart registry locations detects new malware installations.",
  "value": "Registry Run keys are the most common persistence mechanism for malware. Monitoring autostart registry locations detects new malware installations.",
  "implementation": "Sysmon EventCode 13 (RegistryValueSet) monitors registry modifications. Track all autostart locations: Run, RunOnce, RunServices, Winlogon Shell/Userinit, Explorer Shell Folders, and AppInit_DLLs. Filter known-legitimate entries (Program Files, System32). Alert on entries pointing to temp directories, AppData, user profiles, or encoded/obfuscated paths. Monitor both HKLM (system-wide) and HKCU (per-user). MITRE ATT&CK T1547.001.",
  "visualization": "Table (new Run key entries), Alert on suspicious paths, Timeline.",
  "cimModels": [
    "Change"
  ],
  "cimSpl": "| tstats `summariesonly` count\n  from datamodel=Change.All_Changes\n  by All_Changes.user All_Changes.object_category All_Changes.action span=1h\n| sort -count",
  "references": [
    {
      "title": "Splunk_TA_windows",
      "url": "https://splunkbase.splunk.com/app/742"
    },
    {
      "title": "CIM: Change",
      "url": "https://docs.splunk.com/Documentation/CIM/latest/User/Change"
    }
  ],
  "equipment": [
    "windows"
  ],
  "detailedImplementation": "Prerequisites\n• Install and configure the required add-on or app: `Splunk_TA_windows`.\n• Ensure the following data sources are available: `sourcetype=WinEventLog:Microsoft-Windows-Sysmon/Operational` (EventCode 13).\n• For app installation, inputs.conf, and Splunk directory layout, see the Implementation guide: docs/implementation-guide.md\n\nStep 1 — Configure data collection\nSysmon EventCode 13 (RegistryValueSet) monitors registry modifications. Track all autostart locations: Run, RunOnce, RunServices, Winlogon Shell/Userinit, Explorer Shell Folders, and AppInit_DLLs. Filter known-legitimate entries (Program Files, System32). Alert on entries pointing to temp directories, AppData, user profiles, or encoded/obfuscated paths. Monitor both HKLM (system-wide) and HKCU (per-user). MITRE ATT&CK T1547.001.\n\nStep 2 — Create the search and alert\nRun the following SPL in Search (then save as report or alert; adjust time range and threshold as needed):\n\n```spl\nindex=wineventlog EventCode=13\n| where match(TargetObject, \"(?i)(CurrentVersion\\\\\\\\Run|CurrentVersion\\\\\\\\RunOnce|Winlogon\\\\\\\\Shell|Winlogon\\\\\\\\Userinit|Explorer\\\\\\\\Shell Folders)\")\n| where NOT match(Details, \"(?i)(program files|windows\\\\\\\\system32|syswow64)\")\n| table _time, host, User, Image, TargetObject, Details\n| sort -_time\n```\n\nUnderstanding this SPL\n\n**Registry Run Key Persistence Monitoring** — Registry Run keys are the most common persistence mechanism for malware. Monitoring autostart registry locations detects new malware installations.\n\nDocumented **Data sources**: `sourcetype=WinEventLog:Microsoft-Windows-Sysmon/Operational` (EventCode 13). **App/TA** (typical add-on context): `Splunk_TA_windows`. The SPL below should target the same indexes and sourcetypes you configured for that feed—rename `index=` / `sourcetype=` if your deployment differs.\n\nThe first pipeline stage scopes events using **index**: wineventlog.\n\n**Pipeline walkthrough**\n\n• Scopes the data: index=wineventlog. Cross-check against **Data sources** above so indexes and sourcetypes match your ingestion.\n• Filters the current rows with `where match(TargetObject, \"(?i)(CurrentVersion\\\\\\\\Run|CurrentVersion\\\\\\\\RunOnce|Winlogon\\\\\\\\Shell|Winlogon\\\\\\\\Userinit|Expl…` — typically the threshold or rule expression for this monitoring goal.\n• Filters the current rows with `where NOT match(Details, \"(?i)(program files|windows\\\\\\\\system32|syswow64)\")` — typically the threshold or rule expression for this monitoring goal.\n• Pipeline stage (see **Registry Run Key Persistence Monitoring**): table _time, host, User, Image, TargetObject, Details\n• Orders rows with `sort` — combine with `head`/`tail` for top-N patterns.\n\nOptional CIM / accelerated variant (same use case, normalized fields via Common Information Model):\n\n```spl\n| tstats `summariesonly` count\n  from datamodel=Change.All_Changes\n  by All_Changes.user All_Changes.object_category All_Changes.action span=1h\n| sort -count\n```\n\nUnderstanding this CIM / accelerated SPL\n\n**Registry Run Key Persistence Monitoring** — Registry Run keys are the most common persistence mechanism for malware. Monitoring autostart registry locations detects new malware installations.\n\nDocumented **Data sources**: `sourcetype=WinEventLog:Microsoft-Windows-Sysmon/Operational` (EventCode 13). **App/TA** (typical add-on context): `Splunk_TA_windows`. The SPL below should target the same indexes and sourcetypes you configured for that feed—rename `index=` / `sourcetype=` if your deployment differs.\n\nThis **CIM or accelerated** block uses normalized field names and/or `tstats` over data models. Enable **acceleration** on the referenced models (and correct CIM knowledge objects) or the search may return nothing.\n\n**Pipeline walkthrough**\n\n• Uses `tstats` against accelerated summaries for data model `Change.All_Changes` — enable acceleration for that model.\n• Orders rows with `sort` — combine with `head`/`tail` for top-N patterns.\n\nEnable Data Model Acceleration (and metric indexes for `mstats`) for the models or datasets referenced above; otherwise `tstats`/`mstats` may return no results from summaries.\n\n\nStep 3 — Validate\nConfirm that events are present in the index and that the search returns expected results. Compare with known good/bad scenarios if applicable. Verify field extractions and index permissions.\n\nStep 4 — Operationalize\nAdd the search to a dashboard or set up alert actions (email, webhook, PagerDuty, etc.) as required. Document the use case in your runbook and assign an owner. Consider visualizations: Table (new Run key entries), Alert on suspicious paths, Timeline.",
  "grandmaExplanation": "Registry Run keys are the most common persistence mechanism for malware — so you find out before users do when something is slowing down or breaking."
}
