# content/cat-01-server-compute/UC-1.1.122.json
{
  "$schema": "../../schemas/uc.schema.json",
  "id": "1.1.122",
  "title": "Systemd Unit State Monitoring",
  "criticality": "high",
  "difficulty": "intermediate",
  "monitoringType": [
    "Availability",
    "Fault"
  ],
  "splunkPillar": "Observability",
  "dataSources": "`systemctl list-units` output, systemd journal",
  "app": "`Splunk_TA_nix` (scripted input)",
  "spl": "index=os sourcetype=systemd_units host=* NRestarts>0\n| stats sum(NRestarts) as total_restarts by host, Unit\n| where total_restarts > 5\n| sort -total_restarts",
  "cimModels": [
    "N/A"
  ],
  "description": "Track failed/inactive systemd services, auto-restart counts, and service startup time to prevent cascading failures and identify misconfigured or unhealthy units.",
  "value": "Track failed/inactive systemd services, auto-restart counts, and service startup time to prevent cascading failures and identify misconfigured or unhealthy units.",
  "grandmaExplanation": "Track failed/inactive systemd services, auto-restart counts, and service startup time to prevent cascading failures and identify misconfigured or unhealthy units — so you find out before users do when something is slowing down or breaking.",
  "implementation": "Create a scripted input that runs `systemctl list-units --all --no-pager --plain` and `systemctl show --property=ActiveState,SubState,NRestarts` for critical units. Parse ActiveState, SubState, and NRestarts. Run every 60 seconds. For startup time, use `systemd-analyze` output. Alert on any failed units; alert when NRestarts exceeds 5 in 1 hour for critical services.",
  "detailedImplementation": "Prerequisites\n• Install and configure the required add-on or app: `Splunk_TA_nix` (scripted input).\n• Ensure the following data sources are available: `systemctl list-units` output, systemd journal.\n• For app installation, inputs.conf, and Splunk directory layout, see the Implementation guide: docs/implementation-guide.md\n\nStep 1 — Configure data collection\nCreate a scripted input that runs `systemctl list-units --all --no-pager --plain` and `systemctl show --property=ActiveState,SubState,NRestarts` for critical units. Parse ActiveState, SubState, and NRestarts. Run every 60 seconds. For startup time, use `systemd-analyze` output. Alert on any failed units; alert when NRestarts exceeds 5 in 1 hour for critical services.\n\nStep 2 — Create the search and alert\nRun the following SPL in Search (then save as report or alert; adjust time range and threshold as needed):\n\n```spl\nindex=os sourcetype=systemd_units host=* NRestarts>0\n| stats sum(NRestarts) as total_restarts by host, Unit\n| where total_restarts > 5\n| sort -total_restarts\n```\n\nUnderstanding this SPL\n\n**Systemd Unit State Monitoring** — Track failed/inactive systemd services, auto-restart counts, and service startup time to prevent cascading failures and identify misconfigured or unhealthy units.\n\nDocumented **Data sources**: `systemctl list-units` output, systemd journal. **App/TA** (typical add-on context): `Splunk_TA_nix` (scripted input). The SPL below should target the same indexes and sourcetypes you configured for that feed—rename `index=` / `sourcetype=` if your deployment differs.\n\nThe first pipeline stage scopes events using **index**: os; **sourcetype**: systemd_units. If that sourcetype is not mentioned in Data sources, double-check parsing or update the documentation to match the feed you actually ingest.\n\n**Pipeline walkthrough**\n\n• Scopes the data: index=os, sourcetype=systemd_units. Cross-check against **Data sources** above so indexes and sourcetypes match your ingestion.\n• `stats` rolls up events into metrics; results are split **by host, Unit** so each row reflects one combination of those dimensions (useful for per-host, per-user, or per-entity comparisons for this use case).\n• Filters the current rows with `where total_restarts > 5` — typically the threshold or rule expression for this monitoring goal.\n• Orders rows with `sort` — combine with `head`/`tail` for top-N patterns.\n\n\nStep 3 — Validate\nConfirm that events are present in the index and that the search returns expected results. Compare with known good/bad scenarios if applicable. Verify field extractions and index permissions.\n\nStep 4 — Operationalize\nAdd the search to a dashboard or set up alert actions (email, webhook, PagerDuty, etc.) as required. Document the use case in your runbook and assign an owner. Consider visualizations: Table (failed/inactive units by host), Single value (count of failed units), Timechart of restart counts.\n\nScripted input (generic example)\nThis use case relies on a scripted input. In the app's local/inputs.conf add a stanza such as:\n\n```ini\n[script://$SPLUNK_HOME/etc/apps/YourApp/bin/collect.sh]\ninterval = 300\nsourcetype = your_sourcetype\nindex = main\ndisabled = 0\n```\n\nThe script should print one event per line (e.g. key=value). Example minimal script (bash):\n\n```bash\n#!/usr/bin/env bash\n# Output metrics or events, one per line\necho \"metric=value timestamp=$(date +%s)\"\n```\n\nFor full details (paths, scheduling, permissions), see the Implementation guide: docs/implementation-guide.md",
  "visualization": "Table (failed/inactive units by host), Single value (count of failed units), Timechart of restart counts.",
  "references": [
    {
      "title": "Splunk_TA_nix",
      "url": "https://splunkbase.splunk.com/app/833"
    }
  ],
  "equipment": [
    "linux"
  ]
}
